# rootfs/etc/sv/runsvdir-server/run
#!/bin/sh
exec chpst -u server runsvdir /home/server/.sv
